# test_run_at_boot.service
[Unit]
Description=test_run_at_boot

[Service]
ExecStart=/home/ubuntu/misc/service/test_run_at_boot/start.sh

[Install]
WantedBy=multi-user.target

